# conf/vca.cron
# VCA cron
0 * * * * root /usr/bin/php /usr/share/vca/cron/vca.cron.hour.php >> /var/log/vca_cron.log
* * * * * root /usr/bin/php /usr/share/vca/cron/vca.cron.min.php >> /var/log/vca_cron.log
